# root/etc/services.d/dotproxy/run
#!/usr/bin/with-contenv bash

setcap cap_net_bind_service=+ep /app/dotproxy && \
exec \
	s6-setuidgid abc /app/dotproxy -config /config/config.yaml -verbosity ${LOG_LEVEL:-info}